# setup/81-thumbdrive.rules
KERNEL=="sda1" SUBSYSTEM=="block" RUN+="/home/pi/fixtureinspector/setting/usb_restart.sh"
